# backend/agente-suporte.service
/etc/systemd/system/agente-suporte.service 


[Unit]
Description=Agente Suporte NotebookLM - FastAPI
After=network.target

[Service]
Type=simple
Environment=PATH=/home/arpa/agente-suporte-notebooklm/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
WorkingDirectory=/home/arpa/agente-suporte-notebooklm
ExecStart=/home/arpa/agente-suporte-notebooklm/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 75
EnvironmentFile=/home/arpa/agente-suporte-notebooklm/backend/.env
Environment=PYTHONUNBUFFERED=1
User=root
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target